# stage5/06-install-pypilot/files/gpsd.socket
[Unit]
Description=GPS (Global Positioning System) Daemon Sockets

[Socket]
ListenStream=/var/run/gpsd.sock
ListenStream=2947
SocketMode=0600

[Install]
WantedBy=sockets.target